# deploy/systemd/eidocs-prune.service
[Unit]
Description=Prune EI Docs completed/failed job metadata

[Service]
Type=oneshot
WorkingDirectory=/opt/eidocs/current
Environment=EIDOCS_ROOT=/var/lib/eidocs
Environment=EIDOCS_CONFIG_DIR=/etc/eidocs
ExecStart=/opt/eidocs/current/.venv/bin/eidocs prune --older-than-days 7 --storage /var/lib/eidocs
StandardOutput=append:/var/log/eidocs/prune.service.log
StandardError=append:/var/log/eidocs/prune.service.log
Nice=10
NoNewPrivileges=yes
PrivateTmp=yes
MemoryMax=1G
CPUQuota=50%
TasksMax=64
